# home-services/fisher/fisher-update.service
[Unit]
Description=Update fisher plugins

[Service]
Type=oneshot
ExecStart=fisher-update.fish
User=%u
