# home-services/hblock/hblock-update.service
[Unit]
Description=Update hblock list

[Service]
Type=oneshot
ExecStart=hblock
ExecStart=%h/bin/hblock-update.sh
User=%u
