# home-services/hblock/hblock-update.service
[Unit]
Description=Update hblock list

[Service]
Type=oneshot
ExecStart=hblock
User=%u
