# home-services/sdkman/sdkman-update.service
[Unit]
Description=Update sdkman

[Service]
Type=oneshot
ExecStart=%h/bin/sdkman-update.sh
User=%u
