# home-services/nvim/nvim-update.service
[Unit]
Description=Update neovim plugins

[Service]
Type=oneshot
ExecStart=%h/bin/nvim-update.sh
User=%u
