# home-services/nvim/nvim-update.timer
[Unit]
Description=Update neovim plugins
After=network.target

[Timer]
OnCalendar=07:00
OnCalendar=13:00
OnCalendar=19:00
Persistent=true
RandomizedDelaySec=10m
Unit=nvim-update.service

[Install]
WantedBy=timers.target
