# home-services/vs-code/vscode-exts-update.service
[Unit]
Description=Update visual studio code extensions

[Service]
Type=oneshot
ExecStart=code --update-extensions
User=%u
